# contrib/init/migcoind.conf
description "Migcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env MIGCOIND_BIN="/usr/bin/migcoind"
env MIGCOIND_USER="migcoin"
env MIGCOIND_GROUP="migcoin"
env MIGCOIND_PIDDIR="/var/run/migcoind"
# upstart can't handle variables constructed with other variables
env MIGCOIND_PIDFILE="/var/run/migcoind/migcoind.pid"
env MIGCOIND_CONFIGFILE="/etc/migcoin/migcoin.conf"
env MIGCOIND_DATADIR="/var/lib/migcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # migcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$MIGCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run migcoind."
        echo "The setting must appear in $MIGCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Migcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$MIGCOIND_PIDDIR"
    chmod 0755 "$MIGCOIND_PIDDIR"
    chown $MIGCOIND_USER:$MIGCOIND_GROUP "$MIGCOIND_PIDDIR"
    chown $MIGCOIND_USER:$MIGCOIND_GROUP "$MIGCOIND_CONFIGFILE"
    chmod 0660 "$MIGCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$MIGCOIND_PIDFILE" \
    --chuid $MIGCOIND_USER:$MIGCOIND_GROUP \
    --exec "$MIGCOIND_BIN" \
    -- \
    -pid="$MIGCOIND_PIDFILE" \
    -conf="$MIGCOIND_CONFIGFILE" \
    -datadir="$MIGCOIND_DATADIR" \
    -disablewallet \
    -daemon

